# vespai_fixed.service
[Unit]
Description=VespAI Hornet Detection System
After=network.target

[Service]
Type=simple
User=vespai
Group=vespai
WorkingDirectory=/home/vespai
Environment="PATH=/home/vespai/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/vespai/start_vespai.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target